# yoho-remote-server.service
[Unit]
Description=Yoho Remote Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yoho
WorkingDirectory=/var/lib/yoho-remote
EnvironmentFile=/etc/yoho-remote/server.env
ExecStart=/opt/yoho-remote/yoho-remote-server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
