# yoho-remote-daemon.service
[Unit]
Description=Yoho Remote Daemon
After=network.target

[Service]
Type=simple
User=guang
WorkingDirectory=/home/workspaces/repos/yoho-remote
Environment=CLI_API_TOKEN=rDhnX0JCPIki0s6t1kNsHJkSLCvpAEt3wNCb_dkEyOc
Environment=YOHO_REMOTE_URL=http://127.0.0.1:3006
Environment=PATH=/home/guang/.local/bin:/home/guang/.nvm/versions/node/v22.18.0/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/workspaces/repos/yoho-remote/.env
ExecStart=/home/workspaces/repos/yoho-remote/cli/dist-exe/bun-linux-x64/yoho-remote-daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
